# lims_asgi.service
[Unit]
Description=LIMS ASGI Service with Daphne
After=network.target

[Service]
# ########################################################################
# ## 重要：请将下面的 'your_user' 替换为您在服务器上运行此应用的用户名 ##
# ##       例如：User=root 或 User=www-data                          ##
# ########################################################################
User=your_user
Group=your_user

# ########################################################################
# ## 重要：请将下面的 '/path/to/your/project' 替换为您项目的绝对路径    ##
# ##       例如：WorkingDirectory=/home/user/myproject/lims          ##
# ########################################################################
WorkingDirectory=/path/to/your/project

# ########################################################################
# ## 重要：如果您的 daphne 不在 /usr/local/bin/daphne，请修改此路径   ##
# ##       您可以通过在虚拟环境中运行 `which daphne` 来找到它         ##
# ########################################################################
ExecStart=/usr/local/bin/daphne -b 0.0.0.0 -p 8000 lims.asgi:application

Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target

